# services/radio_service.service
[Unit]
Description=UNO Q radio service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arduino
WorkingDirectory=/home/arduino/scripts
ExecStart=/usr/bin/python3 /home/arduino/scripts/radio_service.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target